# build/debian/supervisr.service
[Unit]
Description=supervisr
After=syslog.target network.target
Requires=network.target supervisr-task-runner.service

[Service]
User=supervisr
Group=supervisr
WorkingDirectory=/usr/share/supervisr
Type=simple
ExecStart=/usr/bin/supervisr-ctl run.web

[Install]
WantedBy=multi-user.target

