# build/debian/supervisr-task-runner.service
[Unit]
Description=supervisr Task Runner
After=syslog.target network.target
Requires=network.target

[Service]
User=supervisr
Group=supervisr
WorkingDirectory=/usr/share/supervisr
Type=simple
ExecStart=/usr/bin/sv run.worker

[Install]
WantedBy=multi-user.target

